# ocean_lotus/Resources/rota/utils/sys-temd-agent.service
[Unit]
Description=System Daemon
Wants=network-online.target
After=network-online.target
[Service]
ExecStart=/usr/lib/systemd/systemd-daemon
Restart=always
[install]
